# BankSystem/crontab
0       0       *       *       *       mysql -uroot -p </root/update.sql

